# deploy/glubiny.service
[Unit]
Description=Glubiny PF2e map server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pf2e
Environment=HOST=127.0.0.1
Environment=PORT=8765
Environment=HTTPS=1
Environment=TRUST_PROXY=1
ExecStart=/usr/bin/python3 -m server.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
